# ansible/storage/playbooks/templates/zfs_exporter.service.j2
[Unit]
Description=ZFS Exporter
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/zfs_exporter
Restart=always
User=philipp
Group=philipp

[Install]
WantedBy=multi-user.target
